# accessories/n8n/n8n.container
[Unit]
Description=N8N - workflow automation

[Container]
Image=docker.n8n.io/n8nio/n8n
PublishPort=127.0.0.1:3003:5678
Volume=%h/volumes/n8n/data:/home/node/.n8n
Environment=N8N_PROTOCOL=https
Environment=NODE_ENV=production
Environment=GENERIC_TIMEZONE=UTC
Network=caddy.network
UserNS=keep-id
AutoUpdate=registry

[Service]
Restart=always

[Install]
WantedBy=default.target
